# vbackup.service
[Unit]
Description=Run VBackup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.duplicacy/vbackup
