# vbackup.timer
[Unit]
Description=Run VBackup

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
